# contrib/ddojo.crontab
MAILTO=""
TEST=1
0 0 * * * pi /home/pi/ddojo_local/scripts/update.pi.sh
